# branches/squeeze/xenedu-hpraid/etc/cron.d/xenedu-hpraid
11 */6	* * *	root	/sbin/xenedu-check-hp-raid-courriel

